# rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid mysql /usr/bin/mysqld_safe  > /dev/null 2>&1